# cron.d/desktopchanger_cron
SHELL=/bin/bash
PATH=~/bin:/usr/bin/:/bin
* * * * * <user> /bin/bash -c ". ~/.bashrc; <path_to_project_folder>/notify_test_delete.bash  > <path_to_project_folder>/wpcron.log 2>&1"
